# bin/etc_initd_accumulo
#! /usr/bin/env bash

# Licensed to the Apache Software Foundation (ASF) under one or more
# contributor license agreements.  See the NOTICE file distributed with
# this work for additional information regarding copyright ownership.
# The ASF licenses this file to You under the Apache License, Version 2.0
# (the "License"); you may not use this file except in compliance with
# the License.  You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

#
# accumulo            This shell script takes care of stopping and starting 
#                     accumulo on a single server
# chkconfig: - 60 20
# description: Apache Accumulo: an implementation of Google's BigTable
# probe: false

. /etc/rc.d/init.d/functions

# let's assume this file is symlinked in from accumulo distro bin directory
# otherwise reset the ACCUMULO_HOME setting
pwd=`pwd -P`
ACCUMULO_HOME=`cd $pwd/.. ; pwd -P`

RETVAL=0
case "$1" in
   start)
        $ACCUMULO_HOME/bin/start-here.sh
        RETVAL=$?
   ;;
   stop)
        $ACCUMULO_HOME/bin/stop-here.sh
        RETVAL=$?
   ;;
   status)
        COUNT=`pgrep -f "$ACCUMULO_HOME.*org.apache.accumulo.start" | wc -l`
        echo "There are $COUNT accumulo processes on this machine"
   ;;
   restart)
        $0 stop
        $0 start
   ;;
   *)
      echo "Usage: $0 {start|stop|status|restart}"
      exit 1
   ;;
esac

exit 0
